# bridge-bot/systemd/bridge-bot.service
[Unit]
Description=bridge Telegram control bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=freax
ExecStart=/home/freax/repos/github/freaxnx01/public/bridge/bridge-bot/bridge_bot.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
StandardOutput=append:/home/freax/.cache/bridge/bridge-bot.log
StandardError=inherit

[Install]
WantedBy=multi-user.target
