# clrepo-bot/systemd/clrepo-bot.service
[Unit]
Description=clrepo Telegram wrapper bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/projects/repos/github/freaxnx01/public/clrepo/clrepo-bot/clrepo_bot.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
StandardOutput=append:%h/.cache/clrepo/clrepo-bot.log
StandardError=inherit

[Install]
WantedBy=default.target
